# relay/deploy/yaver-relay.service
[Unit]
Description=Yaver Relay Server — P2P tunnel for mobile-to-desktop connectivity
Documentation=https://github.com/kivanccakmak/yaver.io
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/yaver-relay serve --quic-port=4433 --http-port=8443
Restart=always
RestartSec=5
StartLimitIntervalSec=60
StartLimitBurst=5

# Logging — all output goes to journald
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yaver-relay

# Performance
LimitNOFILE=65535

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true

[Install]
WantedBy=multi-user.target
